# systemd/openclaw-dashboard.service
[Unit]
Description=OpenClaw Dashboard
After=network.target

[Service]
Type=simple
WorkingDirectory=__APP_DIR__
EnvironmentFile=__ENV_FILE__
ExecStart=/usr/bin/env node src/index.js
Restart=on-failure
RestartSec=3

[Install]
WantedBy=default.target
